# bm_sbc_gateway.service
[Unit]
Description=bm_sbc_gateway
After=dev-serial0.device dbus.service
Wants=dev-serial0.device dbus.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# The gateway runs as root so it can write to /usr/local/bin/ during a DFU
# binary swap (link/rename require write permission on the directory).
# If the binary is ever moved to a user-writable location, add User= here.
WorkingDirectory=/usr/local/bin

# DFU rollback guard: runs before each start to detect crash loops and
# restore the previous binary if the new one fails to confirm within
# MAX_ATTEMPTS systemd-managed starts.
# ExecStartPre and ExecStart paths are set to ${CMAKE_INSTALL_PREFIX}/bin
# (default: /usr/local/bin). If borealis_sbc installs to a different prefix,
# update the two paths below accordingly, or configure this file as a
# cmake .in template.
ExecStartPre=/usr/local/bin/dfu_prestart.sh

ExecStart=/usr/local/bin/bm_sbc_gateway --init /etc/bm_sbc/gateway/gateway.toml
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
